# filesystem/etc/profile.d/umask.sh
# Permissões de arquivo para o SAMBA
umask 000
