# filesystem/etc/skel/.config/autostart/volumeicon.desktop
[Desktop Entry]
Type=Application
Exec=volumeicon
Name=Volume Icon
Categories=Utility;GTK
Icon=
OnlyShowIn=XFCE;
